# .docker/docking-worker.conf
[program:docking-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --sleep=3 --tries=1 --timeout=120
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=5
redirect_stderr=true
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stopwaitsecs=3600

[program:docking-php-fpm]
command=php-fpm -F -R
autostart=true
autorestart=true
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
process_name=%(program_name)s_%(process_num)02d

[program:docking-nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
process_name=%(program_name)s_%(process_num)02d
